# gunciorn-www.helloc.conf

start on net-device-up
stop on shutdown

respawn

setuid cyc
chdir /home/cyc/sites/www.helloc.site/blog

exec ../env/bin/gunicorn --bind unix:/tmp/123.56.218.23.socket blogproject.wsgi:application
